# config/ibgateway.service
# Service systemd pour IB Gateway
# Installation : sudo cp config/ibgateway.service /etc/systemd/system/
#                sudo systemctl daemon-reload
#                sudo systemctl enable ibgateway
#                sudo systemctl start ibgateway

[Unit]
Description=IB Gateway - Interactive Brokers API Gateway
After=network.target
Before=trading-bot.service
Wants=network-online.target

[Service]
Type=forking
User=trader
Group=trader
WorkingDirectory=/home/trader/labase-trading-alerts

# Démarrage IB Gateway avec Xvfb et IBC
ExecStart=/home/trader/labase-trading-alerts/scripts/start_ibgateway.sh
ExecStop=/home/trader/labase-trading-alerts/scripts/stop_ibgateway.sh

# Redémarrage automatique si crash (mais pas trop souvent)
Restart=on-failure
RestartSec=30

# Timeout avant SIGKILL
TimeoutStartSec=120
TimeoutStopSec=30

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ibgateway

# Variables d'environnement
Environment="DISPLAY=:1"

[Install]
WantedBy=multi-user.target
